# etc/init.d/buildkite-agent-settings
#!/sbin/openrc-run
depend() {
    after net network-online logger
    before buildkite-agent
    provide buildkite-agent-settings
}

start() {
    cp /etc/buildkite-agent/buildkite-agent.cfg.template /etc/buildkite-agent/buildkite-agent.cfg

    # fetch name
    name=$(curl --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/name" -H "Metadata-Flavor: Google")
    status=$?
    if [ $status -eq 0 ]; then
        echo name="$name-%n" >> /etc/buildkite-agent/buildkite-agent.cfg
    fi

    # fetch token
    token=$(curl --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/attributes/buildkite-token" -H "Metadata-Flavor: Google")
    status=$?
    if [ $status -ne 0 ]; then
        echo "Unable to get buildkite-token"
        exit $status
    fi
    echo token="$token" >> /etc/buildkite-agent/buildkite-agent.cfg

    # fetch priority
    priority=$(curl --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/attributes/buildkite-priority" -H "Metadata-Flavor: Google")
    if [ $? -eq 0 ]; then
        echo priority=$priority >> /etc/buildkite-agent/buildkite-agent.cfg
    fi

    # fetch tags
    tags=$(curl --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/attributes/buildkite-tags" -H "Metadata-Flavor: Google")
    if [ $? -eq 0 ]; then
        echo tags="$tags" >> /etc/buildkite-agent/buildkite-agent.cfg
    fi

    # fetch sshkey
    key=$(curl --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/attributes/buildkite-sshkey" -H "Metadata-Flavor: Google")
    status=$?
    if [ $status -eq 0 ]; then
        mkdir /home/buildkite/.ssh 2>&1 || true
        echo $key | base64 --decode > /home/buildkite/.ssh/identity
        cat <<'END' > /home/buildkite/.ssh/config
    Host *
        IdentityFile /home/buildkite/.ssh/identity
END
        chown -hR buildkite:buildkite /home/buildkite/.ssh
        chmod -R 0700 /home/buildkite/.ssh
    fi

    # fetch instance count
    agent_count=$(curl --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/attributes/buildkite-agent-count" -H "Metadata-Flavor: Google")
    status=$?
    if [ $status -eq 0 ]; then
        for i in $(seq 1 $agent_count); do
            cat <<'END' > /etc/init.d/buildkite-agent-$i
#!/sbin/openrc-run
command="/usr/sbin/buildkite-agent"
command_args="start"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
command_user=buildkite
depend() {
    after net network-online logger
    provide buildkite-agent
}
END
            chmod 0700 /etc/init.d/buildkite-agent-$i
            rc-update add buildkite-agent-$i default
            rc-service buildkite-agent-$i start
        done
    fi

}