# etc/init.d/docker-credential-file
#!/sbin/openrc-run
command_user=buildkite
depend() {
    after net network-online logger docker
    provide docker-credential-file
}

start() {
    while [ ! -S /var/run/docker.sock ]
    do
        echo "Waiting for docker to start..."
        sleep 2
    done
    echo "Docker started..."
    curl --output /tmp/docker-credential-file.encoded --fail --silent "http://metadata.google.internal/computeMetadata/v1/instance/attributes/docker-credential-file" -H "Metadata-Flavor: Google"
    status=$?
    if [ $status -eq 0 ]; then
        cat /tmp/docker-credential-file | base64 -d > /tmp/docker-credential-file
        sudo --user buildkite docker login --username _json_key --password "$(cat /tmp/docker-credential-file)" https://eu.gcr.io
    fi
}
