# etc/init.d/google-instance-setup
#!/sbin/openrc-run
depend() {
    after net network-online logger
    before sshd
    provide google-instance-setup
}

start() {
    /usr/bin/google_instance_setup
}