# etc/init.d/google-accounts-daemon
#!/sbin/openrc-run
command="/usr/bin/google_accounts_daemon"
command_args=""
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true

depend() {
    after net network-online logger google-instance-setup google-network-daemon
    provide google-accounts-daemon
}