# etc/init.d/google-shutdown-scripts
#!/sbin/openrc-run
depend() {
    after net network-online logger google-instance-setup google-network-daemon
    provide google-shutdown-scripts
}

start() {
    /usr/bin/google_metadata_script_runner --script-type shutdown
}